# deploy/systemd/zentriq-worker.service
[Unit]
Description=Zentriq Documents Celery Worker
After=network.target mariadb.service redis.service

[Service]
Type=simple
User=zentriq
Group=zentriq
WorkingDirectory=/opt/zentriq-documents
EnvironmentFile=/opt/zentriq-documents/.env
ExecStart=/opt/zentriq-documents/.venv/bin/celery -A celery_worker.celery worker --loglevel=info
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
