# deploy/systemd/zentriq.service
[Unit]
Description=Zentriq Documents (Gunicorn)
After=network.target mariadb.service redis.service

[Service]
Type=simple
User=zentriq
Group=zentriq
WorkingDirectory=/opt/zentriq-documents
EnvironmentFile=/opt/zentriq-documents/.env
ExecStart=/opt/zentriq-documents/.venv/bin/gunicorn --workers 4 --bind 127.0.0.1:8000 run:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
